# supervisord.conf
[program:autoack]
command=python autoack.py autoack
[supervisord]
logfile=./autoack.log
logfile_maxbytes=5MB
loglevel=debug
